# dist/docker/conf/supervisord/supervisord-cronjob.conf
[supervisord]
nodaemon=true

[program:deploy]
command = /bin/bash -c "/var/www/html/bin/deploy"
autostart = true
autorestart = false
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=nginx
priority=1
exitcodes=0,2

[program:cronjob]
command = /bin/bash -c "/var/www/html/bin/runCronjob"
user = root
autostart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
